# contrib/start_scripts/linux
#!/bin/sh

# chkconfig: 2345 99 01
# description: Bricolage CMS

# Created by David Wheeler, 2005.

# This is an example of a start/stop script for SysV-style init, such as is
# used on Linux systems. You should configure these variables:

# Where to find Bricolage.
BRICOLAGE_ROOT="/usr/local/bricolage"

# Whether to use bric_queued. Set to non-null to make it so.
USE_BRIC_QUEUED=

# Username and password to use with bric_queued
DIST_USERNAME=
DIST_PASSWORD=

# The interval between bric_queued runs.
DIST_INTERVAL=30

# Maximum number of distribution jobs to select for each bric_queued run.
# 0 is no maximum
DIST_CHUNKS=0

# Various file locations.
BRIC_APACHECTL="$BRICOLAGE_ROOT/bin/bric_apachectl"
BRIC_QUEUED="$BRICOLAGE_ROOT/bin/bric_queued"
DIST_PIDFILE="$BRICOLAGE_ROOT/log/bric_queued.pid"
DIST_LOG="$BRICOLAGE_ROOT/log/bric_queued.log"

# Then just copy it to your start scripts directory (usually /etc/init.d) and
# add it to the startup process:

# % cp contrib/start_scripts/linux /etc/init.d/bricolage
# % chmod +x /etc/init.d/bricolage
#
# Then, either create symlinks to the startup script:
#
# % ln -s /etc/init.d/bricolage /etc/rc.d/rc0.d/K02bricolage
# % ln -s /etc/init.d/bricolage /etc/rc.d/rc1.d/K02bricolage
# % ln -s /etc/init.d/bricolage /etc/rc.d/rc2.d/K02bricolage
# % ln -s /etc/init.d/bricolage /etc/rc.d/rc3.d/S98bricolage
# % ln -s /etc/init.d/bricolage /etc/rc.d/rc4.d/S98bricolage
# % ln -s /etc/init.d/bricolage /etc/rc.d/rc5.d/S98bricolage
#
# Or, if you have chkconfig, simply:
#
# % /sbin/chkconfig --add bricolage

##############################################################################

$BRIC_APACHECTL "$@"

if [ $USE_BRIC_QUEUED ]; then
    case "$1" in
        start)
            "$BRIC_QUEUED" --username "$DIST_USERNAME" \
                --password "$DIST_PASSWORD" \
                --log "$DIST_LOG" \
                --pid "$DIST_PIDFILE" --delay "$DIST_INTERVAL"
                --chunks "$DIST_CHUNKS"
            echo bric_queued: started
            ;;
        stop)
            if [ -f "$DIST_PIDFILE" ]; then
                kill `cat "$DIST_PIDFILE"`
                echo bric_queued: stopped
            fi
            ;;
        restart)
            if [ -f "$DIST_PIDFILE" ]; then
                kill `cat "$DIST_PIDFILE"`
            fi
            "$BRIC_QUEUED" --username "$DIST_USERNAME" \
                --password "$DIST_PASSWORD" \
                --log "$DIST_LOG" \
                --pid "$DIST_PIDFILE" --delay "$DIST_INTERVAL"
                --chunks "$DIST_CHUNKS"
            echo bric_queued: restarted
            ;;
    esac;
fi
